# deploy/pvadmin.service
[Unit]
Description=PV Admin Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pvbot
EnvironmentFile=/opt/pvbot/.env
Environment=ADMIN_WEB_ENABLE=1
ExecStart=/opt/pvbot/.venv/bin/python /opt/pvbot/admin_web.py
Restart=always
RestartSec=3
StandardOutput=append:/opt/pvbot/logs/runtime.log
StandardError=append:/opt/pvbot/logs/runtime.log

[Install]
WantedBy=multi-user.target
