# deploy/pvadmin.service
[Unit]
Description=PV Admin Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pvbot/usdt_telegram_membership
EnvironmentFile=/opt/pvbot/usdt_telegram_membership/.env
Environment=ADMIN_WEB_ENABLE=1
ExecStart=/opt/pvbot/usdt_telegram_membership/.venv/bin/python /opt/pvbot/usdt_telegram_membership/admin_web.py
Restart=always
RestartSec=3
StandardOutput=append:/opt/pvbot/usdt_telegram_membership/logs/runtime.log
StandardError=append:/opt/pvbot/usdt_telegram_membership/logs/runtime.log

[Install]
WantedBy=multi-user.target
