# deploy/pvbot-watchdog.timer
[Unit]
Description=PV Bot Watchdog Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=pvbot-watchdog.service

[Install]
WantedBy=timers.target

